# config/crontab
* * * * * cd /var/www/html && setsid /usr/local/bin/php console.php cron > /proc/1/fd/1 2>&1

