# config/crontab
* * * * * www-data /usr/local/bin/php /var/www/html/index.php cron > /proc/1/fd/1 2>&1

